# wsl-x11-vsock/x11-vsock.service
[Unit]
Description=WSL2 VSOCK X11 Forwarding

[Service]
ExecStart=/usr/bin/socat - VSOCK-CONNECT:2:6000
StandardInput=socket
StandardOutput=socket
